# files/fail2ban.upstart
description "fail2ban - ban hosts that cause multiple authentication errors"

start on filesystem and started networking
stop on deconfiguring-networking

expect daemon
respawn

pre-start script
    [ -d /var/run/fail2ban ] || mkdir -p /var/run/fail2ban
end script

exec /usr/bin/fail2ban-client -x start

pre-stop exec /usr/bin/fail2ban-client stop

post-stop exec rm -f /var/run/fail2ban/fail2ban.pid
