# files/fail2ban.upstart
description "fail2ban - ban hosts that cause multiple authentication errors"

start on filesystem and started networking
stop on deconfiguring-networking

pre-start exec /usr/bin/fail2ban-client -x start
post-stop exec /usr/bin/fail2ban-client stop
