# files/fail2ban.upstart
description "fail2ban - ban hosts that cause multiple authentication errors"

start on filesystem and started networking
stop on deconfiguring-networking

expect daemon
respawn

exec /usr/bin/fail2ban-server -x -b

post-stop exec rm -f /var/run/fail2ban/fail2ban.pid
