# files/gentoo-initd
#!/sbin/runscript
# This file is part of Fail2Ban.
#
# Fail2Ban is free software; you can redistribute it and/or modify
# it under the terms of the GNU General Public License as published by
# the Free Software Foundation; either version 2 of the License, or
# (at your option) any later version.
#
# Fail2Ban is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU General Public License
# along with Fail2Ban; if not, write to the Free Software
# Foundation, Inc., 59 Temple Place, Suite 330, Boston, MA  02111-1307  USA
#
# Author: Sireyessire, Cyril Jaquier
# 
# $Revision: 559 $

opts="start stop restart reload showlog"

FAIL2BAN="/usr/bin/fail2ban-client ${FAIL2BAN_OPTIONS}"

depend() {
	need net
	need logger
	after iptables
}

start() {
	ebegin "Starting fail2ban"
	${FAIL2BAN} start &> /dev/null
	eend $? "Failed to start fail2ban"
}

stop() {
	ebegin "Stopping fail2ban"
	${FAIL2BAN} stop &> /dev/null
	eend $? "Failed to stop fail2ban"
}

restart() {
	if ! service_stopped "${SVCNAME}" ; then
		svc_stop || return "$?"
		sleep 1
	fi
	svc_start
}

reload() {
	ebegin "Reloading fail2ban"
	${FAIL2BAN} reload > /dev/null
	eend $? "Failed to reload fail2ban"
}

showlog(){
	less /var/log/fail2ban.log
}
